# templates/4s-boss.service.epp
<%- |
  Stdlib::Absolutepath $path = '/usr/local/bin'
| -%>
[Unit]
Description=4s-boss service
After=network.target

[Service]
Type=simple
ExecStart=<%= $path %>/4s-boss -D
ExecStartPost=/bin/sleep 5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=4s-boss
User=4store
Group=4store
Restart=always

[Install]
WantedBy=multi-user.target
